# systemd/system/plexmediaserver.service.d/override.conf
#
#
[Service]
#  These values are only needed if you wish to change user & group
User=felix
Group=felix

# This is needed to change the default umask   
Umask=0002
